# debian/postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]
then
  invoke-rc.d ferm6 enable || true
  update-rc.d ferm6 start || true
fi

#DEBHELPER#
